# LinuxSoftware_RaspberryPi/home/Raspi_Project/RaspiBot_funV1/code/build/Script
#! /bin/sh
### BEGIN INIT INFO
# Provides:          RaspibotScript
# Required-Start:   
# Required-Stop:     
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: RaspibotScript
# Description:       RaspibotScript
### END INIT INFO

case "$1" in
  start)
	echo "Clean all"
	sudo killall RaspiBot
	
	echo "Running App"
	sudo /home/pi/Raspi_Project/Raspibot/RaspiBot_funV1/code/build/RaspiBot
    ;;
  stop)
	echo "Clean all"
	sudo killall RaspiBot
    ;;

esac
